# deploy/defensor-vial.service
[Unit]
Description=Defensor Vial MX - bot de Telegram (long-polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
